# cron/crontab
* * * * * root sh /root/hello.sh >> /var/log/cron.log 2>&1
# Empty Line